# connect_script.conf
[program:connect_script]
command=/home/boris-zverkov/Neurocar/connect.sh
autostart=false
autorestart=unexpected
stderr_logfile=/var/log/connect.err.log
stdout_logfile=/var/log/connect.out.log

