# script.conf
[program:script]
command=/path/to/script
autostart=false
autorestart=unexpected
stderr_logfile=/path/to/script.err.log
stdout_logfile=/path/to/script.out.log
